# pkgs/by-name/ve/ventoy-compat-protocol/src/run-ventoy.mount
[Mount]
What=/dev/ventoy
Where=/run/ventoy
